# src/supervisord.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf

[program:php-fpm]
command=/usr/sbin/php-fpm7

[supervisord]
logfile=/var/log/supervisord.log 
nodaemon=true
